# templates/tocker@.path
[Path]
PathExists=/run/systemd/transient/tocker_%i.service
Unit=tocker_container_logger@%i.service

[Unit]
Description=watching the transient service tocker_%i
CollectMode=inactive-or-failed
